# debian/minos-tools-extra.backup-mozilla.cron.d
15 */4 * * * root /usr/bin/backup-mozilla
