# debian/minos-core-tools.update-cd.cron.d
#This file is part of the minos-core-tools package
0 23 * * *  root /usr/bin/update-cd
